# wallpaper/.config/systemd/user/wallpaper-hotplug.service
[Unit]
Description=Reload hyprpaper wallpaper on monitor hotplug

[Service]
Type=oneshot
ExecStart=%h/.config/wallpaper/set_wallpaper.sh
